# sys-cluster/apache-zookeeper/files/init
#!/sbin/runscript

export CLASSPATH="$(java-config-2 -d -p apache-zookeeper):${CLASSPATH}";
export ZOOCFGDIR=/etc/zookeeper

start() {
zkServer.sh start > /dev/null
}

stop() {
zkServer.sh stop
}

restart() {
zkServer.sh restart > /dev/null
}

status() {
zkServer.sh status
}
